# profile.d/path.sh
PATH="$GOPATH/bin:$PATH"
#PATH="node_modules/.bin:$PATH"
PATH="$PATH:/opt/local/bin:/opt/local/sbin"
PATH="$PATH:$HOME/src/google/depot_tools"
PATH="$PATH:$HOME/.yarn/bin"
PATH="$PATH:/opt/local/lib/postgresql96/bin"
PATH="$HOME/elements/bin:$PATH"
export PATH;
